# conf/supervisor_logstash.conf
[program:logstash]
command=/opt/logstash-1.4.2/bin/logstash -f /opt/logstash/conf/logstash_simple.conf
